# supervisord.conf
[program:postgres]
command=/etc/init.d/postgresql start
user=postgres
startsecs=0
priority=1
autostart=true

